# Password Sync/Server/Systemd/Dispetcher.service
[Unit]
Description=Password Sync Event Dispatcher
After=network.target

[Service]
Type=simple
User=root
ExecStart=/opt/password_sync/bin/dispatcher.sh
Restart=always

[Install]
WantedBy=multi-user.target
